# hassio-ihost-silabs-multiprotocol/rootfs/etc/s6-overlay/s6-rc.d/zigbeed-socket/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start EmberZNet Zigbee Daemon
# ==============================================================================
# shellcheck disable=SC2015
bashio::log.info "Starting zigbeed-socket..."
exec /usr/bin/socat pty,link=/tmp/ttyZigbee pty,link=/tmp/ttyZigbeeNCP
